# live-build/config/includes.chroot/etc/skel/.config/autostart/kiosk-app.desktop
[Desktop Entry]
Type=Application
Name=Security Monitor Launch
Exec=xfce4-terminal --fullscreen -e btop
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
